# ops/systemd/investor-intel-api.service
[Unit]
Description=Investor Intel API
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/investor-intel
EnvironmentFile=/etc/investor-intel/investor-intel.env
ExecStart=/usr/bin/npm run start:api
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
