# ops/systemd/investor-intel-fast-refresh.service
[Unit]
Description=Investor Intel fast refresh worker
After=network.target postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/opt/investor-intel
EnvironmentFile=/etc/investor-intel/investor-intel.env
ExecStart=/usr/bin/npm run worker:fast
User=www-data
Group=www-data
