# ops/systemd/investor-intel-hourly-ingest.timer
[Unit]
Description=Run Investor Intel hourly ingest

[Timer]
OnCalendar=hourly
Persistent=true
Unit=investor-intel-hourly-ingest.service

[Install]
WantedBy=timers.target
