# ops/systemd/investor-intel-brief.timer
[Unit]
Description=Run Investor Intel brief generation every hour

[Timer]
OnCalendar=hourly
Persistent=true
Unit=investor-intel-brief.service

[Install]
WantedBy=timers.target
